# data/alfred.conf
description "Alfred domotic server"
author "joseph.piron@gmail.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec python -m alfred --db_host hal -d